# scripts/com.flutur.db-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.flutur.db-backup</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/alessioivoycazzaniga/Projects/social-cli-mcp/scripts/backup.sh</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>DATABASE_URL</key>
        <string>postgresql://flutur:dev@localhost:5433/flutur_dev</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>4</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/alessioivoycazzaniga/.config/flutur-ops/backups/launchd-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/alessioivoycazzaniga/.config/flutur-ops/backups/launchd-stderr.log</string>
</dict>
</plist>
